# provisioning/roles/server/templates/vsftpd.service.j2
[Unit]
Description=vsftpd {{ vsftpd_version }} FTP server (lab target)
After=network.target

[Service]
Type=simple
ExecStart={{ vsftpd_binary }} {{ vsftpd_config }}
Restart=on-failure

[Install]
WantedBy=multi-user.target
